# profile.d/02-files.ps1
# File Utilities - Split into Maintainable Fragments
#
# This file has been split into smaller, more maintainable fragments:
#
# - 02-files-navigation.ps1: Directory navigation utilities (.. , ..., ...., ~, desktop, downloads, docs)
# - 02-files-listing.ps1: Enhanced file listing commands (ll, la, lx, tree, bat-cat)
# - 02-files-conversion.ps1: Data and media format conversions (JSON/YAML/CSV/XML/Base64, media, documents)
# - 02-files-utilities.ps1: File inspection and utility functions (head, tail, file-hash, filesize, hex-dump)
#
# Each fragment maintains the lazy loading pattern for optimal startup performance.
