# profile.d/10-wsl.ps1
# ===============================================
# 10-wsl.ps1
# WSL helpers and shorthands
# ===============================================

function wsl-shutdown { wsl --shutdown }
function wsl-list { wsl --list --verbose }
function ubuntu { wsl -D Ubuntu @args }












