# profile.d/99-test-missing-module.ps1
Import-Module NonExistentModule -ErrorAction SilentlyContinue
